# .cogent/skills/detecting-malicious-scheduled-tasks-with-sysmon/forge-meta.json
{
  "name": "detecting-malicious-scheduled-tasks-with-sysmon",
  "description": "Detect malicious scheduled task creation and modification using Sysmon Event IDs 1 (Process Create for schtasks.exe), 11 (File Create for task XML), and Windows Security Event 4698/4702. The analyst correlates task creation with suspicious parent processes, public directory paths, and encoded command arguments to identify persistence and lateral movement via scheduled tasks. Activates for requests involving scheduled task detection, Sysmon persistence hunting, or T1053.005 Scheduled Task/Job analysis.",
  "imported_at": "2026-06-18T01:21:19Z",
  "metadata": {
    "domain": "cybersecurity",
    "subdomain": "threat-hunting",
    "tags": [
      "sysmon",
      "scheduled-tasks",
      "persistence",
      "detection",
      "threat-hunting",
      "windows-security"
    ],
    "version": "1.0",
    "author": "mahipal",
    "license": "Apache-2.0",
    "d3fend_techniques": [
      "Execution Isolation",
      "Process Termination",
      "Hardware-based Process Isolation",
      "Platform Monitoring",
      "Process Suspension"
    ],
    "nist_csf": [
      "DE.CM-01",
      "DE.AE-02",
      "DE.AE-07",
      "ID.RA-05"
    ],
    "mitre_attack": [
      "T1046",
      "T1057",
      "T1082",
      "T1083",
      "T1021"
    ]
  },
  "resources": [
    "scripts/agent.py",
    "references/api-reference.md"
  ]
}