# .cogent/skills/hunting-for-suspicious-scheduled-tasks/forge-meta.json
{
  "name": "hunting-for-suspicious-scheduled-tasks",
  "description": "Hunt for adversary persistence and execution via Windows scheduled tasks by analyzing task creation events, suspicious task properties, and unusual execution patterns that indicate T1053.005 abuse.",
  "imported_at": "2026-06-18T01:21:19Z",
  "metadata": {
    "domain": "cybersecurity",
    "subdomain": "threat-hunting",
    "tags": [
      "threat-hunting",
      "scheduled-tasks",
      "persistence",
      "mitre-t1053-005",
      "windows",
      "endpoint-detection"
    ],
    "version": "1.0",
    "author": "mahipal",
    "license": "Apache-2.0",
    "nist_csf": [
      "DE.CM-01",
      "DE.AE-02",
      "DE.AE-07",
      "ID.RA-05"
    ],
    "mitre_attack": [
      "T1046",
      "T1057",
      "T1082",
      "T1083",
      "T1547"
    ]
  },
  "resources": [
    "scripts/agent.py",
    "scripts/process.py",
    "references/api-reference.md",
    "references/standards.md",
    "references/workflows.md",
    "assets/template.md"
  ]
}